# net/pius/pkg-plist
bin/pius
include/pius/pius_c_api.h
lib/pius/libpius.so
lib/pius/README
@mode 555
etc/rc.d/pius
@mode
@unexec if cmp -s %D/etc/pius-conf.lua.sample %D/etc/pius-conf.lua; then rm -f %D/etc/pius-conf.lua; fi
etc/pius-conf.lua.sample
@exec if [ ! -f %D/etc/pius-conf.lua ] ; then cp -p %D/%F %B/pius-conf.lua; fi
@dirrmtry lib/pius
