# net/pius/pkg-plist
bin/pius
lib/pius/README
etc/pius-conf.lua
@mode 555
etc/rc.d/pius
@mode
@dirrmtry lib/pius
